# docling-streamlit.service
[Unit]
Description=Docling Streamlit Document Q&A Assistant
After=network.target
Wants=network.target

[Service]
Type=simple
User=docling
Group=docling
WorkingDirectory=/opt/docling-app
Environment=PATH=/opt/docling-app/venv/bin
Environment=PYTHONPATH=/opt/docling-app
ExecStart=/opt/docling-app/venv/bin/streamlit run 5-chat.py --server.port 8501 --server.address 0.0.0.0
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal
SyslogIdentifier=docling-streamlit

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/docling-app
ReadWritePaths=/tmp
ReadWritePaths=/var/log

# Resource limits
LimitNOFILE=65536
MemoryLimit=2G
TimeoutStartSec=30
TimeoutStopSec=10

[Install]
WantedBy=multi-user.target